# infra/docker/crontab
# kakeibo-ai crontab (auto-generated from settings.json)
00 06 * * * cd /app && /usr/bin/python main.py --timeframe daily >> /app/logs/cron_execution.log 2>&1
# 空行が必要
