# infra/docker/crontab
# kakeibo-ai crontab
# 毎日 23:50 に自動レビューを実行
50 23 * * * /app/infra/docker/run_cron.sh
# 空行が必要
